# docker/server/crontab
* * * * * echo "1" > /var/log/cron/cron.log 2>&1
